# build-deb/nuage-bridge/debian/DEBIAN/prerm
#!/bin/bash

/etc/init.d/nuage stop
update-rc.d -f nuage remove
exit 0

