# build-deb/nuage-bridge/debian/DEBIAN/prerm
#!/bin/bash

if [ -f /etc/inittab ] ; then
	/etc/init.d/nuage stop
	update-rc.d -f nuage remove
else
	systemctl stop nuage
	systemctl disable nuage
fi
exit 0
